# selim-start.conf
start on runlevel [2345]
stop on runlevel [!2345]

exec /home/pi/Documents/Selim.io/selim_init.py